# llx-qjack-pulse-integrator.install-files/etc/systemd/user/pulseaudio.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/pulseaudio --realtime=yes --daemonize=no --log-target=journal
